# pulseblog-deploy/pulseblog.service
[Unit]
Description=PulseBlog API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/www/wwwroot/www.kangban.online
ExecStart=/www/server/pulseblog/pulseblog-api
Restart=always
RestartSec=5
Environment=GIN_MODE=release
Environment=SERVER_PORT=8080
Environment=DB_HOST=127.0.0.1
Environment=DB_PORT=3306
Environment=DB_USER=root
Environment=DB_PASSWORD=你的数据库密码
Environment=DB_NAME=blog_system
Environment=JWT_SECRET=你的32位随机字符串
Environment=DEFAULT_ADMIN_USERNAME=admin
Environment=DEFAULT_ADMIN_PASSWORD=你的管理员密码
Environment=WEB_ORIGIN=https://www.kangban.online
Environment=ADMIN_ORIGIN=https://www.kangban.online
Environment=UPLOAD_DIR=./uploads
Environment=DEEPSEEK_API_KEY=你的DeepSeekKey（不用则留空）

[Install]
WantedBy=multi-user.target
